# deploy/supervisor.conf
[program:nurcrm-worker-default]
command=php /var/www/nurcrm/artisan queue:work redis --queue=default --sleep=3 --tries=3 --max-time=3600
directory=/var/www/nurcrm
user=www-data
numprocs=3
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/nurcrm-worker.log
stderr_logfile=/var/log/nurcrm-worker-error.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3

[program:nurcrm-worker-notifications]
command=php /var/www/nurcrm/artisan queue:work redis --queue=notifications --sleep=3 --tries=3 --max-time=3600
directory=/var/www/nurcrm
user=www-data
numprocs=2
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/nurcrm-notifications.log
stderr_logfile=/var/log/nurcrm-notifications-error.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3

[program:nurcrm-horizon]
command=php /var/www/nurcrm/artisan horizon
directory=/var/www/nurcrm
user=www-data
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/nurcrm-horizon.log
stderr_logfile=/var/log/nurcrm-horizon-error.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3

[program:nurcrm-scheduler]
command=php /var/www/nurcrm/artisan schedule:work
directory=/var/www/nurcrm
user=www-data
autostart=true
autorestart=true
stdout_logfile=/var/log/nurcrm-scheduler.log
stderr_logfile=/var/log/nurcrm-scheduler-error.log

[group:nurcrm]
programs=nurcrm-worker-default,nurcrm-worker-notifications,nurcrm-horizon,nurcrm-scheduler
